# roles/zygote/templates/fcgiwrap-zygote.service.j2
[Unit]
Description=Zygote CGI Server
After=nss-user-lookup.target

[Service]
ExecStart=/usr/sbin/fcgiwrap
User={{ zygote_config.user }}
