# roles/zygote/templates/fcgiwrap-zygote.socket.j2
[Unit]
Description=Zygote fcgiwrap socket

[Socket]
ListenStream=/run/fcgiwrap-zygote.socket

[Install]
WantedBy=sockets.target
